# systemd/starter.path
[Unit]
Description=Path activation for the "starter" programme

[Path]
PathExists=/tmp/amd_ready
