# systemd/starter.socket
[Socket]
ListenStream=/tmp/phlock
SocketMode=0777

[Install]
WantedBy=sockets.target
