# packaging/starter.service
[Unit]
Description=Start the "starter" program
Requires=starter-pre.service
After=starter-pre.service wm_ready.service ac.service deviced.service launchpad-preload.service

[Service]
ExecStart=/usr/bin/starter
Restart=always
RestartSec=0
KillSignal=SIGKILL
User=system
Group=system
SmackProcessLabel=starter

[Install]
WantedBy=tizen-system.target
