# packaging/starter-pre.service
[Unit]
Description=Pre actions for "starter"

[Service]
Type=oneshot
ExecStart=/etc/init.d/rd3starter
SmackProcessLabel=starter
