# packaging/starter.service
[Unit]
Description=Start the "starter" program
Before=core-efl.target
After=xorg.target

[Service]
ExecStart=/usr/bin/starter

[Install]
WantedBy=core-efl.target
WantedBy=tizen-mobile-session.target
